# supporting_files/supervisord-openssh.conf
[program:openssh]
command=/start-openssh.sh
numprocs=1
autostart=true
autorestart=true
